# systemd/proxymanager-postgresql.container
[Unit]
Description=Proxy Manager PostgreSQL Database
After=network-online.target

[Container]
Image=docker.io/library/postgres:18-alpine
ContainerName=proxymanager-postgresql
EnvironmentFile=.env
Volume=proxymanager-database.volume:/var/lib/postgresql/data
Pod=proxymanager.pod

# Health check (runs inside container)
HealthCmd=pg_isready -d $${POSTGRES_DB} -U $${POSTGRES_USER}
HealthInterval=30s
HealthTimeout=5s
HealthRetries=5
HealthStartPeriod=20s

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
